# supervisor.d/mqtt.ini
[program:mqtt]
command=mosquitto
redirect_stderr=true
stdout_logfile=storage/logs/mqtt.log
stopsignal=sigint
stopwaitsecs=10

