# supervisor.d/webrtc.ini
[program:webrtc]
command=/usr/bin/php artisan webrtc:listen
redirect_stderr=true
stdout_logfile=storage/logs/webrtc.log
stopsignal=sigint
stopwaitsecs=10

