# supervisor.d/websocket.ini
[program:soketi]
command=npm run websocket
stopsignal=sigint
stopwaitsecs=60
redirect_stderr=true
stdout_logfile=storage/logs/soketi-supervisor.log
